# ansible/roles/app-batch/files/etc/cron.d/jincreek-deauthentication
PATH=/sbin:/bin:/usr/sbin:/usr/bin:/opt/apps/bin
MAILTO=""
10 0 * * * jincreek /opt/apps/bin/batch.sh deauthentication -organization_code 4 > /var/log/jincreek/cron.log 2>&1
10 1-23 * * * jincreek /opt/apps/bin/batch.sh deauthentication -organization_code 4 >> /var/log/jincreek/cron.log 2>&1
